# com.agentic.intel-collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentic.intel-collector</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/SSDRAID0/agentic-system/mcp-servers/world-intel-mcp/.venv/bin/python</string>
        <string>-m</string>
        <string>world_intel_mcp.collector</string>
        <string>--daemon</string>
        <string>--interval</string>
        <string>300</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Volumes/SSDRAID0/agentic-system/mcp-servers/world-intel-mcp</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/tmp/intel-collector.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/intel-collector-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Volumes/SSDRAID0/agentic-system/mcp-servers/world-intel-mcp/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Volumes/SSDRAID0/agentic-system/mcp-servers/world-intel-mcp/src</string>
        <key>WORLD_INTEL_LOG_LEVEL</key>
        <string>INFO</string>
    </dict>
</dict>
</plist>
